# systemd/hots-api.service
[Unit]
Description=HOTS Hot Search API Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/hots-api
ExecStart=/opt/hots-api/hots-api
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=hots-api
Environment="GIN_MODE=release"

[Install]
WantedBy=multi-user.target